# systemd/dotclaw.service
[Unit]
Description=DotClaw Telegram Assistant
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=dotclaw
WorkingDirectory=/home/dotclaw/dotclaw
Environment=NODE_ENV=production
EnvironmentFile=/home/dotclaw/dotclaw/.env
ExecStart=/usr/bin/node /home/dotclaw/dotclaw/dist/index.js
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillMode=process
StandardOutput=append:/home/dotclaw/dotclaw/logs/dotclaw.log
StandardError=append:/home/dotclaw/dotclaw/logs/dotclaw.error.log

[Install]
WantedBy=multi-user.target
